# .codex/skills/hooyodex-service-restart/assets/systemd/hooyodex.service
[Unit]
Description=Hooyodex Telegram Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/path/to/hooyodex --config %h/.config/hooyodex/config.toml run
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3
TimeoutStopSec=20
Environment=RUST_LOG=hooyodex=info
Environment=PATH=/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:%h/.local/bin
WorkingDirectory=%h

[Install]
WantedBy=default.target
