# deploy/codexclaw.service
[Unit]
Description=CodexClaw Telegram Orchestrator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=__BIN__ --config %h/.config/codexclaw/config.toml
WorkingDirectory=%h/.config/codexclaw
Environment=CODEXCLAW_CONFIG=%h/.config/codexclaw/config.toml
Environment=RUST_LOG=info
Restart=on-failure
RestartSec=5
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=default.target
